# rootfs/etc/services.d/a1mini/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

: "${APP_DATA_DIR:=/data}"
# 确保数据目录的工作文件夹存在
mkdir -p "${APP_DATA_DIR}/3mf" "${APP_DATA_DIR}/temp" "${APP_DATA_DIR}/tempDownload"

: "${APP_CODE_DIR:="${APP_DATA_DIR}/addon_code"}"

# 默认后备代码目录（镜像内只读）
DEFAULT_CODE_DIR="/usr/src/app"

# 合并候选目录，顺序：显式 APP_CODE_DIR -> addon_configs 匹配目录 -> 常规位置
echo "临时调试-CANDIDATES: ${CANDIDATES}"
# 优先使用第一个名字包含 bambu_queue_print 或 bambu_queue 的子目录

# 检查：/config/bambu_queue_print
TARGET_DIR="/config/bambu_queue_print"

CODE_DIR="${DEFAULT_CODE_DIR}"
if [ -d "${TARGET_DIR}" ] && [ -f "${TARGET_DIR}/app.py" ]; then
  CODE_DIR="${TARGET_DIR}"
  echo "调试模式-使用覆盖路径运行: ${CODE_DIR}"
else
  echo "未发现覆盖文件，使用默认路径: ${CODE_DIR}"
fi

# 强制容器内监听地址和端口，确保 Supervisor ingress 能连接到容器
# Supervisor 会把 ingress 转发到容器的 5000 端口，确保绑定 0.0.0.0
export HOST="0.0.0.0"
export PORT="5000"

echo "Starting app with HOST=${HOST} PORT=${PORT} APP_DATA_DIR=${APP_DATA_DIR} CODE_DIR=${CODE_DIR}"
ps aux || true

# 使用镜像内 venv 的 python 运行选定目录下的 app.py（这样可以复用已安装依赖）
exec /usr/src/app/venv/bin/python -u "${CODE_DIR}/app.py"
